# deploy/single-node/mc-control-plane-reconciler.service
[Unit]
Description=mc-control-plane durable Operation reconciler
After=network-online.target mc-control-plane-host-api.service
Wants=network-online.target
PartOf=mc-control-plane.target

[Service]
Type=simple
User=opc
Group=opc
EnvironmentFile=/etc/mc-control-plane/secrets.env
ExecStart=/usr/local/bin/mc-control-plane node-reconciler
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
NoNewPrivileges=yes
PrivateTmp=yes
UMask=0077

[Install]
WantedBy=mc-control-plane.target
